# bridge/discord-dm/discord-dm-bridge.service
# systemd unit for the Discord DM bridge on the OCI box.
# Install: sudo cp discord-dm-bridge.service /etc/systemd/system/
#          sudo systemctl daemon-reload && sudo systemctl enable --now discord-dm-bridge
# Logs:    journalctl -u discord-dm-bridge -f
#
# Runs from /opt (NOT /home): on SELinux-enforcing Oracle Linux, systemd
# services can't read EnvironmentFile/WorkingDirectory under /home — it fails
# with "Permission denied". Deploy to /opt/discord-dm and `restorecon -R` it
# (see README). User `opc` is the Oracle Linux default; swap to `ubuntu` on an
# Ubuntu image. `node` is the official tarball at /usr/local/bin/node.
#
# StartLimit* live in [Unit] (systemd ignores them in [Service]).

[Unit]
Description=MilesVault Discord DM bridge
After=network-online.target
Wants=network-online.target
# Back off if it crash-loops (Discord rate-limits aggressive reconnects).
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=opc
WorkingDirectory=/opt/discord-dm
EnvironmentFile=/opt/discord-dm/.env
ExecStart=/usr/local/bin/node index.mjs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
